# roles/keystone/templates/keystone-public.service.j2
[Unit]
Description=Keystone-admin
After=network.target

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/bin/keystone-wsgi-public --port 5000

[Install]
WantedBy=multi-user.target